# templates/jupyterhub.service.epp
# Managed by Puppet - manual changes will be overwritten.
[Unit]
Description=Jupyterhub

[Service]
ExecStart=<%= $jupyterhub::pyvenv %>/bin/jupyterhub --config=<%= $jupyterhub::config_dir %>/jupyterhub_config.py
WorkingDirectory=<%= $jupyterhub::jupyterhub_dir %>
Environment=PATH=<%= $jupyterhub::pyvenv %>/bin:$PATH
User=<%= $jupyterhub::jupyterhub_username %>

[Install]
WantedBy=multi-user.target
